# launchd/io.local.mac-bootstrap.claude-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  DEPRECATED 2026-06-12
  Replaced by: the tmux-based ai-work workspace launcher in scripts/tmux-workspace.sh.
  Reason: tmux now owns the daemon window and ai-work pane layout directly.
  Do not re-enable this launch agent.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.local.mac-bootstrap.claude-daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>{{BOOTSTRAP}}/scripts/claude-daemon-tmux.sh</string>
	</array>
	<key>KeepAlive</key>
	<false/>
	<key>RunAtLoad</key>
	<false/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>StandardOutPath</key>
	<string>/tmp/claude-daemon-tmux.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/claude-daemon-tmux.err</string>
	<key>Disabled</key>
	<true/>
</dict>
</plist>
