# launchd/io.local.mac-bootstrap.claude-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.local.mac-bootstrap.claude-daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>{{BOOTSTRAP}}/scripts/claude-daemon-tmux.sh</string>
	</array>
	<key>KeepAlive</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>StandardOutPath</key>
	<string>/tmp/claude-daemon-tmux.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/claude-daemon-tmux.err</string>
</dict>
</plist>
