# launchd/io.local.mac-bootstrap.claude-keepalive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  DEPRECATED 2026-06-11
  Replaced by: claude-daemon-tmux.sh with StartCalendarInterval (00:00 / 08:00 / 15:00).
  Reason: per-minute polling (StartInterval 60) deemed unnecessary.
  Do NOT install or re-enable this launch agent.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.local.mac-bootstrap.claude-keepalive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>{{BOOTSTRAP}}/scripts/claude-daemon-keepalive.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>60</integer>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardOutPath</key>
	<string>/tmp/claude-daemon-keepalive.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/claude-daemon-keepalive.err</string>
</dict>
</plist>
