# launchd/io.local.mac-bootstrap.devspace.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.local.mac-bootstrap.devspace</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>{{BOOTSTRAP}}/scripts/devspace-supervisor.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>{{REPO_ROOT}}</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>{{LOG_DIR}}/launchd-devspace.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>{{LOG_DIR}}/launchd-devspace.stderr.log</string>
</dict>
</plist>
